# deploy/supervisor/orvae-queue-mail.conf
; Worker solo para la cola "emails" (mismo QUEUE_CONNECTION=database).
; 1) En Laravel, los envíos en cola deben usar ->onQueue('emails') o $queue = 'emails' en el Mailable/Job.
; 2) Copiar a /etc/supervisor/conf.d/orvae-queue-mail.conf
; 3) sudo supervisorctl reread && sudo supervisorctl update && sudo supervisorctl start orvae-queue-mail:*
;
[program:orvae-queue-mail]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/orvaepe/artisan queue:work database --queue=emails --sleep=3 --tries=3 --max-time=3600
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
numprocs=1
redirect_stderr=true
stdout_logfile=/var/www/orvaepe/storage/logs/queue-mail.log
stopwaitsecs=3600
