# deploy/supervisor/orvae-queue-worker.conf
; Copiar a /etc/supervisor/conf.d/orvae-queue-worker.conf
; Ajusta: command (ruta al proyecto), user, numprocs si quieres varios workers.
;
; Tras instalar:
;   sudo supervisorctl reread
;   sudo supervisorctl update
;   sudo supervisorctl start orvae-queue-worker:*
;
; Tras cada deploy con cambios en código o config:
;   php artisan queue:restart
;
[program:orvae-queue-worker]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/orvaepe/artisan queue:work database --sleep=3 --tries=3 --max-time=3600
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
numprocs=1
redirect_stderr=true
stdout_logfile=/var/www/orvaepe/storage/logs/queue-worker.log
stopwaitsecs=3600
